# radon4/etc/systemd/system/rabbitmq-server.service.d/90-limits.conf
# This file managed by Puppet - DO NOT EDIT
[Service]
LimitNOFILE=16384
